# packages/canopycms-cdk/worker/canopy-worker.service
[Unit]
Description=CanopyCMS Worker Daemon
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/canopy-worker
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=/opt/canopy-worker/.env

[Install]
WantedBy=multi-user.target
